# linux/fleet/api.service
[Unit]
Description=Vajra Fleet API service.

[Service]
Type=simple
ExecStart=/bin/bash /home/sablearjun/vajra-fleet/cmake-build-debug/api.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
